# .config/systemd/user/acpid.service
[Unit]
Description=%N service

[Service]
ExecStart=/usr/bin/acpid -nfS -c %h/.config/acpi/events

[Install]
WantedBy=default.target
