# .config/systemd/user/apod.service
[Unit]
Description=%N service

[Service]
Type=oneshot
SyslogIdentifier=%N
ExecStart=/bin/bash -c 'apod get'
ExecStart=-/bin/bash -c 'apod set'
