# .config/systemd/user/daily.timer
[Unit]
Description=daily timer

[Timer]
OnCalendar=8:00
Persistent=yes

[Install]
WantedBy=timers.target
