# .config/systemd/user/env-update.service
[Unit]
Description=%N service
DefaultDependencies=no
Before=sockets.target timers.target paths.target basic.target

[Service]
Type=oneshot
SyslogIdentifier=%N
RemainAfterExit=yes
Environment=DISPLAY=:1
ExecStart=/bin/sh -lc 'unset SHLVL; systemctl --user import-environment'
