# .config/systemd/user/jackdbus.service
[Unit]
Description=%N service

[Service]
ExecStart=/usr/bin/jackdbus auto
ExecStop=/usr/bin/jack_control exit
