# .config/systemd/user/mbsync@.service
[Unit]
Description=%N service

[Service]
Type=oneshot
SyslogIdentifier=%N
ExecStart=/bin/sh -c 'mbsync-bin %i'
