# .config/systemd/user/mbsync@head.timer
[Unit]
Description=%N timer

[Timer]
OnActiveSec=0
OnUnitInactiveSec=10m
