# .config/systemd/user/mbsync@.service
[Unit]
Description=%N service

[Service]
SyslogIdentifier=%N
ExecStart=/usr/bin/bash -c 'mbsync-bin %i'
