# .config/systemd/user/mbsync@tail.timer
[Unit]
Description=%N timer

[Timer]
OnCalendar=*:00,30
Persistent=yes
