# .config/systemd/user/mbsync@tail.timer
[Unit]
Description=mbsync(tail) timer

[Timer]
OnCalendar=*:00,30
Persistent=yes
