# .config/systemd/user/mbsync@head.timer
[Unit]
Description=mbsync(head) timer

[Timer]
OnActiveSec=0
OnUnitInactiveSec=5m
