# .config/systemd/user/mbsync@tail.timer
[Unit]
Description=mbsync(tail) timer

[Timer]
OnCalendar=*:30:00
Persistent=yes
