# .config/systemd/user/mbsync.timer
[Unit]
Description=mbsync timer

[Timer]
OnUnitInactiveSec=10m
