# .config/systemd/user/mopidy.service
[Unit]
Description=mopidy

[Service]
Environment=DISPLAY=:1
ExecStart=/usr/bin/mopidy

[Install]
WantedBy=default.target
