# .config/systemd/user/ssh-tunnel@.service
[Unit]
Description=ssh-tunnel(%i) service

[Service]
Type=simple
SyslogIdentifier=ssh-tunnel(%i)
ExecStart=/bin/bash -lc 'autossh -M 0 -F ~/.ssh/config.tunnels -N "%i"'
Restart=on-failure

[Install]
WantedBy=default.target
