# .config/systemd/user/xsession@.target
[Unit]
Description=%N target

Wants=urxvtd@%i.service
