# .termux/sv/ssh-tunnel@baobab/run
#!/usr/bin/env bash

srv=${PWD##*/}
host=${srv##*@}

exec autossh -M 0 -F ~/.ssh/config.tunnels -N "$host"
